# com.sellermind.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sellermind.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/MAMP/bin/php/php8.4.1/bin/php</string>
        <string>/Applications/MAMP/htdocs/sellermind-ai/artisan</string>
        <string>schedule:work</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Applications/MAMP/htdocs/sellermind-ai</string>

    <key>StandardOutPath</key>
    <string>/Applications/MAMP/htdocs/sellermind-ai/storage/logs/scheduler-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Applications/MAMP/htdocs/sellermind-ai/storage/logs/scheduler-stderr.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Applications/MAMP/bin/php/php8.4.1/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
