# supervisor-sellermind-worker.conf
[program:sellermind-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/sellermind-ai/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/path/to/sellermind-ai/storage/logs/worker.log
stopwaitsecs=3600
priority=999

; Переменные окружения
environment=LANG=ru_RU.UTF-8,LC_ALL=ru_RU.UTF-8
